# templates/fairdata-etsin-qvain.socket.j2
[Unit]
Description=UNIX socket for Etsin-Qvain

[Socket]
ListenStream={{ gunicorn_socket }}
SocketMode=0666
SocketUser={{ application_user }}
SocketGroup={{ application_group }}

[Install]
WantedBy=sockets.target
